# systemd/cot-translator.service
[Unit]
Description=CoT Translator (GPGGA UDP -> CoT -> TAK TLS)
After=network-online.target
Wants=network-online.target

[Service]
User=cot
Group=cot
WorkingDirectory=/opt/cot
EnvironmentFile=-/etc/default/cot-translator
ExecStart=/opt/cot/cot-translator
Restart=always
RestartSec=1
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
MemoryMax=150M
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
